# contact-phase-summary-sync.timer
[Unit]
Description=Contact Phase Summary Sync Timer (Deprecated - Use batch-job-queue-scheduler)
Requires=contact-phase-summary-sync.service

[Timer]
OnCalendar=Mon *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target




