# contact-sales-badge.service
[Unit]
Description=HubSpot Contact Sales Badge Update
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/mirai-api
Environment=PATH=/var/www/mirai-api/venv/bin:/usr/bin:/usr/local/bin
ExecStart=/var/www/mirai-api/venv/bin/python3 /var/www/mirai-api/scripts/queue_batch_job.py contact-sales-badge
StandardOutput=journal
StandardError=journal
SyslogIdentifier=contact-sales-badge

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/mirai-api/logs
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target




